# 7.0/supervisord.conf
[supervisord]
nodaemon=true

[program:mysqld]
command=/etc/init.d/mysqld start
autostart=true

[program:httpd]
command=/etc/init.d/httpd start
autostart=true
