# 5.4/supervisord.conf
[supervisord]
nodaemon=true
logfile=/var/log/supervisor/supervisord.log
pidfile=/var/run/supervisord.pid
childlogdir=/var/log/supervisor

[program:mysqld]
command=/etc/init.d/mysqld start
autostart=true

[program:httpd]
command=/etc/init.d/httpd start
autostart=true
